# systemd/hush-control.service
[Unit]
Description=hush-control — fleet console control plane (LAN mode)
Documentation=https://github.com/clarkbar-sys/hush
Conflicts=hush-control-tsnet.service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=hush
Group=hush
WorkingDirectory=/etc/hush
# The :-defaults keep hush-control bootable if control.env is missing (the -
# prefix already makes the file optional): an unset variable would otherwise
# expand to an empty argument and override the binary's own flag default.
EnvironmentFile=-/etc/hush/control.env
ExecStart=/usr/local/bin/hush-control -listen ${HUSH_CONTROL_LISTEN:-127.0.0.1:8080} -config ${HUSH_CONTROL_CONFIG:-/etc/hush/fleet.json}
Restart=on-failure
RestartSec=5s

# Sandboxing. hush-control talks HTTP to agents and persists fleet.json (the
# web console can add/remove machines), so /etc/hush needs to stay writable
# even under ProtectSystem=strict.
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/etc/hush
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
LockPersonality=true

[Install]
WantedBy=multi-user.target
