# systemd/hush-control-tsnet.service
[Unit]
Description=hush-control — fleet console control plane (tsnet mode)
Documentation=https://github.com/clarkbar-sys/hush
Conflicts=hush-control.service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=hush
Group=hush
WorkingDirectory=/etc/hush
EnvironmentFile=-/etc/hush/control-tsnet.env
ExecStart=/usr/local/bin/hush-control -tsnet -hostname ${HUSH_CONTROL_HOSTNAME} -state-dir ${HUSH_CONTROL_STATE_DIR}
Restart=on-failure
RestartSec=5s
StateDirectory=hush

# tsnet serves HTTPS on :443; grant just enough capability to bind it as an
# unprivileged user instead of running as root.
NoNewPrivileges=true
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
ProtectSystem=strict
# hush-control persists fleet.json here when machines are added/removed from
# the web console, so /etc/hush needs to stay writable despite strict mode.
ReadWritePaths=/etc/hush
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
# AF_NETLINK is required: tsnet reads the host's routes/addresses over a netlink
# socket at startup (route ip+net / netlinkrib). Without it the node fails to
# come up with "address family not supported by protocol". This mirrors the
# upstream tailscaled.service allowlist.
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
LockPersonality=true

[Install]
WantedBy=multi-user.target
