# systemd/hush-agent.service
[Unit]
Description=hush-agent — fleet vitals agent
Documentation=https://github.com/clarkbar-sys/hush
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=hush
Group=hush
EnvironmentFile=-/etc/hush/agent.env
ExecStart=/usr/local/bin/hush-agent -listen ${HUSH_AGENT_LISTEN}
Restart=on-failure
RestartSec=5s

# Sandboxing. The agent is read-only and writes nothing, so the whole
# filesystem stays read-only (ProtectSystem=strict) and privilege escalation is
# off. What it can *read* is deliberately left to the Unix layer: the /browse
# endpoint lists directories bounded only by what the unprivileged "hush" user
# can see, so we do NOT set ProtectHome/PrivateTmp here — those would blank out
# /home, /root, /run and the real /tmp regardless of file permissions, making
# the security boundary the sandbox instead of the user's identity. Tighten or
# widen the agent's reach by changing the hush user's group membership
# (usermod -aG ...), not by editing this unit. The kernel-surface protections
# below stay on — they don't hide readable files, they just lock down knobs the
# agent never touches.
NoNewPrivileges=true
ProtectSystem=strict
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
# AF_NETLINK is required for -listen tailnet: resolving the machine's tailnet
# address calls net.InterfaceAddrs(), which enumerates interfaces over a netlink
# socket. Without it the call fails with "address family not supported by
# protocol" and the agent loops "no tailnet address after N attempts" even when
# tailscale0 is up. This mirrors the hush-control-tsnet unit's allowlist.
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
LockPersonality=true

[Install]
WantedBy=multi-user.target
