# systemd/hush-control.service
[Unit]
Description=hush-control — fleet console control plane (LAN mode)
Documentation=https://github.com/clarkbar-sys/hush
Conflicts=hush-control-tsnet.service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=hush
Group=hush
WorkingDirectory=/etc/hush
# Environment= supplies the defaults; the optional EnvironmentFile= below (- prefix)
# overrides them when control.env is present. This keeps hush-control bootable if
# the file is missing — an unset variable would otherwise expand to an empty
# argument and override the binary's own flag default.
#
# NOTE: systemd does NOT support shell-style ${VAR:-default} expansion — that form
# is passed through literally and reaches the binary verbatim, breaking startup.
# The Environment=/EnvironmentFile= pair is the correct way to express a default.
Environment=HUSH_CONTROL_LISTEN=127.0.0.1:8080
Environment=HUSH_CONTROL_CONFIG=/etc/hush/fleet.json
EnvironmentFile=-/etc/hush/control.env
ExecStart=/usr/local/bin/hush-control -listen ${HUSH_CONTROL_LISTEN} -config ${HUSH_CONTROL_CONFIG}
Restart=on-failure
RestartSec=5s

# Sandboxing. hush-control talks HTTP to agents and persists fleet.json (the
# web console can add/remove machines), so /etc/hush needs to stay writable
# even under ProtectSystem=strict.
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/etc/hush
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
LockPersonality=true

[Install]
WantedBy=multi-user.target
