# systemd/hush-agent.service
[Unit]
Description=hush-agent — fleet vitals agent
Documentation=https://github.com/clarkbar-sys/hush
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=hush
Group=hush
# Environment= supplies the default; the optional EnvironmentFile= below (- prefix)
# overrides it when /etc/hush/agent.env is present. If that file is missing — e.g.
# an immutable-distro OS update wiped /etc/hush (SteamOS resets /etc on updates) —
# HUSH_AGENT_LISTEN keeps this default instead of expanding to an empty -listen
# argument (which the agent would read as Go's net default of :80, fail with
# "permission denied" as the unprivileged hush user, and crash-loop).
#
# NOTE: systemd does NOT support shell-style ${VAR:-default} expansion — that form
# is passed through literally and reaches the binary verbatim, breaking startup.
# The Environment=/EnvironmentFile= pair is the correct way to express a default.
Environment=HUSH_AGENT_LISTEN=tailnet
EnvironmentFile=-/etc/hush/agent.env
ExecStart=/usr/local/bin/hush-agent -listen ${HUSH_AGENT_LISTEN}
Restart=on-failure
RestartSec=5s

# Persist the Job scheduler's jobs.json here. StateDirectory creates
# /var/lib/hush as 0700 hush:hush and keeps it writable even under
# ProtectSystem=strict, so no ReadWritePaths is needed and the agent resolves it
# automatically via $STATE_DIRECTORY. It costs nothing when jobs are off (the
# agent only touches it once enabled), so it's set unconditionally.
StateDirectory=hush

# Sandboxing. The agent's only writable spot is its StateDirectory above; the
# rest of the filesystem stays read-only (ProtectSystem=strict) and privilege
# escalation is off. What it can *read* is deliberately left to the Unix layer:
# the /browse endpoint lists directories bounded only by what the unprivileged
# "hush" user can see, so we do NOT set ProtectHome/PrivateTmp here — those would
# blank out /home, /root, /run and the real /tmp regardless of file permissions,
# making the security boundary the sandbox instead of the user's identity.
# Tighten or widen the agent's reach by changing the hush user's group
# membership (usermod -aG ...), not by editing this unit. Note /exec and /jobs
# run commands as that same hush user — both opt-in (jobs is off by default) —
# so they inherit exactly this reach. The kernel-surface protections below stay
# on — they don't hide readable files, they just lock down knobs the agent's own
# code never touches.
NoNewPrivileges=true
ProtectSystem=strict
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
# AF_NETLINK is required for -listen tailnet: resolving the machine's tailnet
# address calls net.InterfaceAddrs(), which enumerates interfaces over a netlink
# socket. Without it the call fails with "address family not supported by
# protocol" and the agent loops "no tailnet address after N attempts" even when
# tailscale0 is up. This mirrors the hush-control-tsnet unit's allowlist.
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
LockPersonality=true

[Install]
WantedBy=multi-user.target
